# deploy/launchd/com.newthinker.atlas.analysis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.newthinker.atlas.analysis</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/zuowei/workspace/runtime/atlas/scripts/ops/trigger-analysis.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/zuowei/workspace/runtime/atlas</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/go/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <!-- 每 30 分钟触发一轮分析。改频率：调整 StartInterval（秒）后 deploy + install-services。
       router 的 4h 冷却防止同一信号重复推送，故 30min 不会刷屏。 -->
  <key>StartInterval</key>
  <integer>1800</integer>

  <!-- 加载时不立即跑（避免装服务即触发一轮通知）。 -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/zuowei/workspace/runtime/atlas/logs/analysis.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/zuowei/workspace/runtime/atlas/logs/analysis.err.log</string>
</dict>
</plist>
